# systemd-system/fstrim.timer.d/override.conf
# Copy the parent of this file i.e. 'fstrim.timer.d' directory into '/etc/systemd/system'.

[Timer]
OnCalendar=daily
RandomizedDelaySec=1min
AccuracySec=5min
Persistent=true
